# newsdedup.service.default
[Unit]
Description=News dedup for Tiny Tiny RSS.

[Service]
Environment=VIRTUAL_ENV=/home/<username>/<path>/newsdedup/.env
Environment=PATH=/home/<username>/<path>/src/newsdedup/.env/bin
ExecStart=/home/<username>/<path>/src/newsdedup/newsdedup.py -D -v -q /home/<username>/<path>/src/newsdedup/newsdedup.cfg
Restart=always
RestartSec=10s
TimeoutSec=0

[Install]
WantedBy=default.target
